# root/etc/s6-overlay/s6-rc.d/init-buildx-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

docker buildx create --driver docker-container --name container >/dev/null 2>&1
